# tenzir/services/systemd/tenzir-node.service.in
[Unit]
Description=Tenzir Node
StartLimitIntervalSec=120
StartLimitBurst=2
Wants=network-online.target
After=network-online.target

[Service]
Type=notify
Restart=always

# user + privileges
# We explicitly don't use DynamicUser because that can lead to recursive `chown`s.
# Doing that is pretty slow on some file systems (e.g., xfs).
User=tenzir
Group=tenzir
NoNewPrivileges=yes

# Directories
StateDirectory=tenzir
CacheDirectory=tenzir
LogsDirectory=tenzir

# capabilities
RestrictNamespaces=yes
RestrictAddressFamilies=
CapabilityBoundingSet=
AmbientCapabilities=CAP_NET_BIND_SERVICE
RestrictSUIDSGID=yes

# limits
LimitNOFILE=65535

# system access
ExecPaths=@CMAKE_INSTALL_FULL_BINDIR@ @CMAKE_INSTALL_FULL_LIBEXECDIR@
ProtectSystem=strict
ReadWritePaths=/tmp
PrivateTmp=no
# Allow read access to the home directory of the user so config files can be
# read. In theory the config directory may be outside of the home directory, but
# systemd has no mechanism to specify that.
ProtectHome=read-only
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes

SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

# service specifics
TimeoutStopSec=600
WorkingDirectory=/var/lib/tenzir
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/tenzir-node

[Install]
Alias=tenzir.service
WantedBy=multi-user.target
